# launchd/com.memoant.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.memoant</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/kylenorthup/Code/memoant/.venv/bin/python</string>
        <string>/Users/kylenorthup/Code/memoant/scripts/watch.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/kylenorthup/.memoant/memoant.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/kylenorthup/.memoant/memoant.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/kylenorthup/Code/memoant</string>
</dict>
</plist>
